# le9uo_patches/supplemental-assets/post-boot/rc.local
#!/bin/sh
# This script is supposed to be replaced or appended to /etc/rc.local to activate le9 at boot-time.

echo n | tee /sys/kernel/mm/lru_gen/enabled

